# services/lidarr/target/soularr.timer
[Unit]
Description=Soularr poll — every 10 minutes

[Timer]
# 10min cadence is the soularr-upstream default; aligns with Lidarr's
# typical wanted-list update frequency. Persistent so a missed fire
# (LXC stopped) runs on next boot.
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
